# deploy/systemd/ccframe-web.service
[Unit]
Description=CCFrame Web (Next.js)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/ccframe
Environment=NODE_ENV=production
EnvironmentFile=/opt/ccframe/.env
ExecStart=/usr/bin/npm start --silent
Restart=always
RestartSec=3
User=www-data

[Install]
WantedBy=multi-user.target

